# files/etc/cron.d/nauta-login
*/1 * * * * /usr/sbin/nauta-login.sh >> /tmp/nauta.log 2>&1
